# launchd/com.kjlee.discord-thread-titler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kjlee.discord-thread-titler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kjlee/Projects/discord-thread-titler/venv/bin/python3</string>
        <string>/Users/kjlee/Projects/discord-thread-titler/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kjlee/Projects/discord-thread-titler</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kjlee/Projects/discord-thread-titler/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kjlee/Projects/discord-thread-titler/logs/bot.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
